# package/LTEFix-custom/customizelte/files/etc/crontabs/root
0,10,20,30,40,50 * * * * /usr/lib/custom/johns_ping.sh